# koreader-highlight-joplin-server/root/etc/services.d/joplin/run
#! /usr/bin/with-contenv bash

# 启动joplin
exec s6-setuidgid koreader joplin server start
