# tt-rss-plugins/root/etc/services.d/caddy/run
#! /usr/bin/with-contenv bash

#启动 caddy
exec caddy run --config /etc/caddy/caddy.conf --adapter caddyfile
