# tt-rss-plugins/root/etc/services.d/caddy/run
#! /usr/bin/with-contenv bash

#启动 caddy
exec s6-setuidgid ttrss caddy  --conf /etc/caddy/caddy.conf 
