# aria2/aria2-shared/root/etc/services.d/cron/run
#! /usr/bin/with-contenv bash

#启动trackers更新任务
crond -f -L /dev/null
