# kms/root/etc/services.d/kms/run
#! /usr/bin/with-contenv bash

# 启动kms
exec  kms vlmcsd -D
