# baidupcs-web/root/etc/services.d/baidupcs-web/run
#! /usr/bin/with-contenv bash

# 启动BaiduPCS-Go
exec BaiduPCS-Go
