# tt-rss-plugins/root/etc/services.d/update_daemon/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid ttrss php8 /usr/local/tt-rss/app/update_daemon2.php > /dev/null 2>&1
