# kms/root/etc/services.d/vlmcsd/run
#! /usr/bin/with-contenv bash

# 启动vlmcsd
exec  vlmcsd -D
