# kms/root/etc/services.d/vlmcsd/run
#! /usr/bin/with-contenv bash

#启动vlmcsd
exec s6-setuidgid kms vlmcsd -D -P $VLMCSD_SERVER_PORT
