# anki-sync-server/root/etc/services.d/nginx/run
#! /usr/bin/with-contenv bash

# 启动nginx
exec nginx -g "daemon off;"
