# tt-rss-plugins/root/etc/services.d/opencc/run
#! /usr/bin/with-contenv bash

#启动 opencc api
if [ $ENABLE_OPENCC == "true" ]; then
	exec s6-setuidgid ttrss node /usr/local/OpenCC.henry.wang/node_modules/opencc.henry.wang/bin/www
esle
	exit 0
fi
