# calibre-web/root/etc/services.d/calibre-web/run
#! /usr/bin/with-contenv bash

# 启动calibre-web
exec  /usr/local/calibre-web/cps.py
