# calibre-web/root/etc/services.d/calibre-web/run
#! /usr/bin/with-contenv bash

# 启动calibre-web
exec s6-setuidgid calibre python3 /usr/local/calibre-web/app/cps.py
