# anki-sync-server/root/etc/services.d/anki-sync-server/run
#! /usr/bin/with-contenv bash

# 启动anki-sync-server
exec s6-setuidgid ankisync python3 -m ankisyncd
