# qiandao/root/etc/services.d/qiandao/run
#! /usr/bin/with-contenv bash

# 启动qiandao程序
exec /usr/local/qiandao/run.py

