# qiandao/root/etc/services.d/qiandao/run
#! /usr/bin/with-contenv bash

# 启动qiandao程序
exec s6-setuidgid qiandao /usr/local/qiandao/run.py
