# aria2/root/etc/services.d/darkhttpd/run
#! /usr/bin/with-contenv bash

#启动darkhttpd
exec s6-setuidgid aria2 darkhttpd /usr/local/aria2/AriaNg --port 8080
