# aria2/root/etc/services.d/darkhttpd/run
#! /usr/bin/with-contenv bash

# 启动darkhttpd
exec darkhttpd /usr/local/aria2/AriaNg
