# koreader-sync-server/root/etc/services.d/koreader-sync-server/run
#! /usr/bin/with-contenv bash

# 启动koreader-sync-server
cd /app/koreader-sync-server
exec s6-setuidgid koreader gin start
