# koreader-highlight-joplin-server/root/etc/services.d/inotify/run
#! /usr/bin/with-contenv bash

# 启动inotify
exec s6-setuidgid koreader bash /usr/local/joplin/inotify.sh
