# tt-rss-plugins/root/etc/services.d/php-fpm/run
#! /usr/bin/with-contenv bash

#启动 php-fpm
exec s6-setuidgid ttrss php-fpm8 -F
