# kms/root/etc/services.d/darkhttpd/run
#! /usr/bin/with-contenv bash

#启动darkhttpd
if [ "$KMS_README_WEB" == "true" ]; then
  exec s6-setuidgid kms darkhttpd /home/www --port $KMS_README_WEB_PORT
else
  s6-svc -Od .
  exit 0
fi
