# kms/root/etc/services.d/darkhttpd/run
#! /usr/bin/with-contenv bash

#启动darkhttpd
if [ "$KMS_README_WEB" == "true" ]; then
  exec darkhttpd /home/www
fi
