# tt-rss/root/etc/services.d/php-fpm/run
#! /usr/bin/with-contenv bash

#启动 php-fpm
s6-setuidgid ttrss  php-fpm7 -F
